# supervisor.conf
[program:rrrewear]
command=node app.js
directory=/src
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/err.log
stdout_logfile=/var/log/supervisor/out.log
environment=MYSQL_HOST="db",DB_PORT="3306",MYSQL_USER="root",MYSQL_PASS="fakepassword",MYSQL_DATABASE="rrrewear" 